# wordpress/crontab
* * * * * cd /var/www/html && php wp-cron.php